# deploy/macos/studio.peacockery.apple-mail-mcp-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>studio.peacockery.apple-mail-mcp-helper</string>
  <key>ProgramArguments</key>
  <array>
    <string>__APPLE_MAIL_MCP_HELPER_EXECUTABLE__</string>
    <string>--serve</string>
    <string>__APPLE_MAIL_MCP_HELPER_SOCKET__</string>
  </array>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>/tmp/apple-mail-mcp-helper.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/apple-mail-mcp-helper.err.log</string>
</dict>
</plist>
